# facebook-media-ocr.service
[Unit]
Description=Facebook Media OCR - Group photo scraper (weekly)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/ryzen/projects/FacebookMediaOcr/run_all.sh --systemd
WorkingDirectory=/home/ryzen/projects/FacebookMediaOcr
User=ryzen
StandardOutput=journal
StandardError=journal
